# deploy/oracle/systemd/strategy-arena-backup.timer
[Unit]
Description=Run Strategy Arena backup daily

[Timer]
OnCalendar=*-*-* 21:30:00
Persistent=true

[Install]
WantedBy=timers.target
